# tools/pr-notify/pr-notify.timer
[Unit]
Description=OSAC PR Slack Notifier - weekday mornings

[Timer]
OnCalendar=Mon..Fri *-*-* 09:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
