# tools/pr-notify/pr-notify.service
# NOTE: Paths below assume the repo is cloned at ~/work/src/github/osac-workspace.
# Edit ExecStart, WorkingDirectory, and --config to match your local checkout.

[Unit]
Description=OSAC PR Slack Notifier - post dashboard summary to Slack
After=network.target

[Service]
Type=oneshot
TimeoutStartSec=60
ExecStart=/usr/bin/python3 %h/work/src/github/osac-workspace/tools/pr-notify/notify.py --config %h/work/src/github/osac-workspace/tools/pr-notify/config.toml
WorkingDirectory=%h/work/src/github/osac-workspace/tools/pr-notify
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin

[Install]
WantedBy=default.target
